# docker/php/supervisor.d/queue-general-worker.ini
[program:queue-worker-general]
command=php /var/www/%(ENV_PROJECT_CODE)s/artisan queue:work --tries=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
